# ref/chistotyumen-laravel-main/docker/supervisor/conf.d/queue.conf.example
; Rename to queue.conf and rebuild the image to run queue workers in the same container.
; Prefer a dedicated worker container or Horizon in production at scale.

[program:queue]
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
user=www
autostart=true
autorestart=true
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
